# services/auto_shutdown.service
[Unit]
Description=Auto Shutdown Service

[Service]
User=ubuntu
WorkingDirectory=/auto_shutdown

ExecStart=/bin/bash /auto_shutdown/auto_shutdown.sh
