# scripts/minecraft.service
[Unit]
Description=Minecraft Service
After=default.target

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/bin/java -jar /minecraft/run.sh
